# singbox-router-webui.service
[Unit]
Description=Sing-box WebUI (API + React UI)
After=network.target
Wants=network.target

[Service]
Type=simple
User=sbweb
Group=sbweb
WorkingDirectory=/opt/sb-webui

# Bind only on localhost; nginx proxies it
Environment=PORT=3001
Environment=BIND=127.0.0.1
Environment=NODE_ENV=production

# Build React UI on every start/restart
ExecStartPre=/usr/bin/npm --prefix /opt/sb-webui/web ci
ExecStartPre=/usr/bin/npm --prefix /opt/sb-webui/web run build

# Start backend
ExecStart=/usr/bin/node /opt/sb-webui/server.js

Restart=on-failure
RestartSec=2
TimeoutStartSec=300

# (optional) a bit of hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target