# singbox-router-webui.service
[Unit]
Description=Sing-box WebUI
After=network.target
Wants=network.target

[Service]
Type=simple
User=sbweb
Group=sbweb
WorkingDirectory=/opt/sb-webui

Environment=PORT=3001
Environment=BIND=127.0.0.1
Environment=NODE_ENV=production
Environment=HOME=/home/sbweb
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# ensure UI exists (fail fast with clear message)
ExecStartPre=/usr/bin/test -f /opt/sb-webui/web/dist/index.html

ExecStart=/usr/bin/node /opt/sb-webui/server.js

Restart=on-failure
RestartSec=2
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target
